# services/backupator-verification.service
[Unit]
Description=Backupator Verification Agent
Wants=basic.target
After=basic.target network.target

[Service]
WorkingDirectory=/opt/backupator
ExecStart=/opt/backupator/bin/backup_verify.sh
KillMode=process
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
